# app-admin/Ocsinventory-Agent/files/ocsinventory-agent.cron
#! /bin/sh

/usr/bin/ocsinventory-agent
