# dev-lang/tcron/files/tcron
#!/sbin/runscript
#
# tcrond          Start/Stop the tcron clock daemon.
#
# chkconfig: 2345 21 60
# description: tcron is a UNIX program that runs user-specified \
#              programs at periodic scheduled times; and it can set rtc
#              to power up/down your PC with ATX support.
# processname: tcrond
# config: /etc/tcrontab
# pidfile: /var/run/tcrond.pid

depend() {
	use clock logger
	need localmount
#	provide cron	
}

start() {
	ebegin "Starting tcron"
	  start-stop-daemon --start --quiet --exec /usr/sbin/tcrond --pidfile /var/run/tcron.pid
	eend $?
}

stop() {
	ebegin "Stopping tcron"
	  start-stop-daemon --stop --quiet /usr/sbin/tcrond --pidfile /var/run/tcron.pid
	eend $?
}
