# net-news/newsd/files/newsd.init
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

start() {
	ebegin "Starting newsd"
	start-stop-daemon --start --quiet --exec /usr/sbin/newsd -- -c /etc/newsd.conf
	eend $?
}

stop() {
	ebegin "Stopping newsd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/newsd
	eend $?
}
