# sys-mobilephone/gsm0710muxd/files/gsm0710muxd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need dbus
}

GSMMUXD_BINARY=${GSMMUXD_BINARY:-/usr/sbin/gsm0710muxd}
GSMMUXD_PIDFILE=${GSMMUXD_PIDFILE:-/var/run/${SVCNAME}.pid}

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --background --make-pidfile \
			  --pidfile ${GSMMUXD_PIDFILE} \
			  --exec ${GSMMUXD_BINARY} \
			  -- -x /sys/bus/platform/devices/neo1973-pm-gsm.0 -f 100
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --pidfile ${GSMMUXD_PIDFILE}
	eend $?
}
