# dev-python/translate-toolkit/files/init
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header:

opts="${opts} refreshstats"
PID_FILE="/var/run/pootle/pid"

depend () {
	need localmount
	after bootmisc
	after net
}

start () {
	ebegin "Starting Pootle"
	start-stop-daemon --start \
		--exec /usr/bin/PootleServer \
		--chuid=pootle:pootle \
		--make-pidfile \
		--pidfile="$PID_FILE" \
		--background \
		-- ${POOTLE_OPTS} \
		>> /var/log/pootle 2>&1
	eend $?
}

stop () {
	ebegin "Stopping Pootle"
	start-stop-daemon --stop \
		--pidfile "$PID_FILE"
	eend $?
}

refreshstats () {
	ebegin "Refreshing Pootle statistics"
	/usr/bin/PootleServer --refreshstats > /dev/null 2>&1
	eend $?
}
