# sys-apps/inetutils/files/rlogind.service
[Unit]
Description=remote login server

[Service]
ExecStart=/usr/libexec/rlogind

