# sys-apps/inetutils/files/rlogind.socket
[Unit]
Description=remote login server

[Socket]
ListenStream=513
#BindIPv6Only=ipv6-only
#BindToDevice=tun0
#MaxConnections=8