# dev-db/couchdb/files/couchdb.init-2.1.0
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

pidfile=${COUCHDB_PID_FILE}
command=${EXEC:-/opt/couchdb/bin/couchdb}
command_args="-o /dev/null -e /dev/null ${COUCHDB_OPTIONS}"
start_stop_daemon_args="--background --make-pidfile --user ${COUCHDB_USER}"

depend() {
	need net
}

start_pre() {
	checkpath -q -d -m 0755 -o ${COUCHDB_USER} /run/couchdb
}

stop_post() {
	killall -u ${COUCHDB_USER%:*} epmd
}
