# net-misc/dhcp-fwd/files/dhcp-fwd.initd
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $


depend() {
	need net
}

start() {
	ebegin "Starting dhcp-fwd"
	export PIDFILE=${PIDFILE}
	/usr/sbin/dhcp-fwd -c ${CONFIG}
	eend $?
}

stop() {
	ebegin "Stopping dhcp-fwd"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
	eend $?
}
