# net-misc/vernemq/files/vernemq.initd
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	use dns
	need epmd net
}

start() {
	ebegin "Starting vernemq"
	/usr/bin/vernemq ${VERNEMQ_OPTS} start
	eend $?
}

stop() {
	ebegin "Stopping vernmq"
	/usr/bin/vernemq stop >/dev/null 2>&1
	eend 0
}
