# mail-mta/citadel/files/citadel.init
#!/sbin/runscript

# original Author: Wilfried Goesgens <citadel@outgesourced.org>
# adapted for Gentoo: the_mgt <themgt@mail.ru>

depend() {
        use dns logger
	need net
        provide mta
}

RUNDIR="/var/run/citadel"
DAEMON="/usr/sbin/citserver"
PIDFILE="${RUNDIR}/citadel.pid"
SENDCOMMAND="/usr/sbin/sendcommand"

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

svc_start() {
	ebegin "Starting citserver"
	start-stop-daemon --start --exec /usr/sbin/citserver -- ${CITSERVER_OPTS}
	eend $?
}

svc_stop() {
	ebegin "Stopping citserver"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
	eend $?
}

restart() {
	ebegin "Restarting citserver"
	svc_stop && sleep 3 && svc_start
	eend $? "Failed to restart Citadel"
}
